# 80autousb.rules
KERNEL=="sd*", SUBSYSTEM=="block", RUN+="/etc/udev/client.sh $env{ACTION} %k"